# roles/container_streaming/templates/ytdl-sub/crontab
# min   hour    day     month   weekday command
  0     */6     *       *       *       /config/cron_cmd.sh > /config/ytdl-sub.log 2>&1
